# com.akshay.propertybot.rentals.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.akshay.propertybot.rentals</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/akshay/Documents/GitHub/real_estate_bot/scheduler.py</string>
        <string>--rentals</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/akshay/Documents/GitHub/real_estate_bot</string>

    <!-- Rental scan: 6pm daily -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>18</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/akshay/Documents/GitHub/real_estate_bot/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/akshay/Documents/GitHub/real_estate_bot/launchd.err</string>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
